# trogdor/auv-systemd-units/auv-uptimed.service
[Unit]
Description=Manages auv-uptimed
#Wants=dev-ttyUSB_thor_
#After=dev-ttyUSB_thor_

[Service]
User=software
ExecStart=/home/software/cuauv/software/trogdor/auv-systemd-set-environment.sh auv-uptimed

[Install]
WantedBy=trogdor.target
