# install/runit/trogdor/run
#!/bin/sh
su -c "exec /home/software/cuauv/workspaces/worktrees/master/trogdor/trogdor_boot.sh" software
